# root/etc/s6-overlay/s6-rc.d/init-mod-jellyfin-psp-ffmpeg/run
#!/usr/bin/with-contenv bash
# Docker Mod pour Jellyfin PSP - Wrapper FFmpeg H.264 optimisé

echo "==================================================="
echo "PSP Jellyfin Mod: Installation wrapper FFmpeg"
echo "==================================================="

FFMPEG_ORIGINAL="/usr/lib/jellyfin-ffmpeg/ffmpeg.original"
FFMPEG_WRAPPER="/usr/lib/jellyfin-ffmpeg/ffmpeg"

# Vérifier si le wrapper est déjà installé
if [ -f "$FFMPEG_ORIGINAL" ]; then
    echo "→ Wrapper FFmpeg déjà installé ✅"
    exit 0
fi

# Renommer l'original FFmpeg
echo "→ Sauvegarde de l'original FFmpeg..."
mv "$FFMPEG_WRAPPER" "$FFMPEG_ORIGINAL"

# Créer le wrapper FFmpeg
echo "→ Création du wrapper FFmpeg PSP H.264..."
cat > "$FFMPEG_WRAPPER" << 'FFMPEG_WRAPPER'
#!/bin/bash
# Wrapper FFmpeg pour PSP - H.264 optimisé
# H.264 AVC Base Profile pour Media Engine PSP
# Décodage hardware via sceMpeg, performances optimales

LOGFILE="/tmp/ffmpeg-psp-wrapper.log"

# Détecter H.264
if echo "$@" | grep -q "h264\|libx264"; then
    echo "========================================" >> "$LOGFILE"
    echo "[$(date)] FFMPEG H.264 DÉTECTÉ" >> "$LOGFILE"
    echo "COMMANDE ORIGINALE:" >> "$LOGFILE"
    printf '%s\n' "$@" >> "$LOGFILE"
    echo "" >> "$LOGFILE"
    
    # Parser et optimiser pour PSP H.264
    ARGS=()
    SKIP_NEXT=0
    INJECT_AFTER_CODEC=0
    
    i=0
    while [ $i -lt $# ]; do
        i=$((i + 1))
        arg="${!i}"
        
        if [ $SKIP_NEXT -eq 1 ]; then
            SKIP_NEXT=0
            continue
        fi
        
        # Supprimer les options conflictuelles avec notre profil H.264
        if [[ "$arg" =~ ^-(preset|crf|force_key_frames)$ ]] || [[ "$arg" = "-b:v" ]] || [[ "$arg" = "-maxrate" ]] || [[ "$arg" = "-bufsize" ]]; then
            echo "  [REMOVE] $arg (conflit avec profil H.264 PSP)" >> "$LOGFILE"
            SKIP_NEXT=1
            continue
        fi
        
        # Ajouter l'argument
        ARGS+=("$arg")
        
        # Injecter nos options H.264 APRÈS -c:v libx264
        if [[ "$arg" = "libx264" ]] && [ $INJECT_AFTER_CODEC -eq 0 ]; then
            echo "  [INJECT] Options H.264 Base Profile PSP" >> "$LOGFILE"
            
            # ===== H.264 BASE PROFILE PSP =====
            # Profil et niveau (Base Profile obligatoire pour PSP)
            ARGS+=("-profile:v" "baseline")   # Base Profile (66)
            ARGS+=("-level" "3.0")            # Level 3.0
            
            # Bitrate et résolution
            ARGS+=("-b:v" "768k")             # Bitrate vidéo
            ARGS+=("-maxrate" "768k")         # Bitrate max
            ARGS+=("-bufsize" "1536k")        # Buffer size
            
            # Résolution PSP native
            ARGS+=("-s" "480x272")            # Résolution PSP
            
            # Framerate
            ARGS+=("-r" "24")                 # 24 fps
            
            # GOP et keyframes
            ARGS+=("-g" "48")                 # GOP size (2 secondes à 24fps)
            ARGS+=("-keyint_min" "48")        # Min GOP
            ARGS+=("-sc_threshold" "0")       # Désactiver scene cut detection
            
            # Pas de B-frames (Base Profile)
            ARGS+=("-bf" "0")
            ARGS+=("-refs" "1")               # 1 reference frame
            
            # Encodage CAVLC (pas CABAC)
            ARGS+=("-coder" "0")              # CAVLC
            
            # Motion estimation
            ARGS+=("-me_method" "hex")        # Motion estimation
            ARGS+=("-subq" "5")               # Subpixel quality
            
            # Pas de partitions complexes
            ARGS+=("-partitions" "none")
            ARGS+=("-trellis" "0")
            
            # Deblocking filter
            ARGS+=("-flags" "+loop")
            ARGS+=("-cmp" "chroma")
            
            # Pas de 8x8 DCT (Base Profile)
            ARGS+=("-8x8dct" "0")
            ARGS+=("-fast_pskip" "1")
            
            # Pas de weighted prediction
            ARGS+=("-wpredp" "0")
            ARGS+=("-mixed_refs" "0")
            
            # Fast start pour MP4
            ARGS+=("-movflags" "+faststart")
            
            INJECT_AFTER_CODEC=1
        fi
    done
    
    echo "COMMANDE H.264 MODIFIÉE:" >> "$LOGFILE"
    printf '%s ' "${ARGS[@]}" >> "$LOGFILE"
    echo "" >> "$LOGFILE"
    echo "" >> "$LOGFILE"
    
    exec /usr/lib/jellyfin-ffmpeg/ffmpeg.original "${ARGS[@]}"
fi

# Si ce n'est pas H.264, passer tel quel
exec /usr/lib/jellyfin-ffmpeg/ffmpeg.original "$@"
FFMPEG_WRAPPER

# Rendre le wrapper exécutable
chmod +x "$FFMPEG_WRAPPER"

echo "→ Wrapper FFmpeg MJPEG installé ✅"
echo "→ Tous les transcodages MJPEG seront optimisés pour PSP !"
echo "==================================================="

