# root/etc/s6-overlay/s6-rc.d/init-mod-jellyfin-psp-ffmpeg/run
#!/usr/bin/with-contenv bash
# Docker Mod pour Jellyfin PSP - Wrapper FFmpeg Baseline Profile

echo "==================================================="
echo "PSP Jellyfin Mod: Installation wrapper FFmpeg"
echo "==================================================="

FFMPEG_ORIGINAL="/usr/lib/jellyfin-ffmpeg/ffmpeg.original"
FFMPEG_WRAPPER="/usr/lib/jellyfin-ffmpeg/ffmpeg"

# Vérifier si le wrapper est déjà installé
if [ -f "$FFMPEG_ORIGINAL" ]; then
    echo "→ Wrapper FFmpeg déjà installé ✅"
    exit 0
fi

# Renommer l'original FFmpeg
echo "→ Sauvegarde de l'original FFmpeg..."
mv "$FFMPEG_WRAPPER" "$FFMPEG_ORIGINAL"

# Créer le wrapper FFmpeg
echo "→ Création du wrapper FFmpeg PSP..."
cat > "$FFMPEG_WRAPPER" << 'FFMPEG_WRAPPER'
#!/bin/bash
# Wrapper FFmpeg pour forcer H.264 Baseline Profile (PSP)
# Parse et modifie les arguments pour injecter les options PSP

# Si ce n'est pas un transcodage H.264, passer tel quel
if ! echo "$@" | grep -q "libx264"; then
    exec /usr/lib/jellyfin-ffmpeg/ffmpeg.original "$@"
fi

# Parser les arguments et injecter nos options PSP
ARGS=()
INJECT_DONE=0

for arg in "$@"; do
    ARGS+=("$arg")
    
    # Après "-codec:v:0 libx264", injecter nos options PSP
    if [ "$arg" = "libx264" ] && [ $INJECT_DONE -eq 0 ]; then
        # Options PSP critiques pour Baseline Profile
        ARGS+=("-bf" "0")                    # Pas de B-frames
        ARGS+=("-refs" "1")                  # 1 seule référence
        ARGS+=("-coder" "0")                 # CAVLC (pas CABAC)
        ARGS+=("-g" "48")                    # GOP fixe 48 frames
        ARGS+=("-sc_threshold" "0")          # Pas de scene detection
        ARGS+=("-keyint_min" "48")           # Keyframe min = GOP
        ARGS+=("-partitions" "none")         # Désactiver partitions avancées
        ARGS+=("-me_method" "hex")           # Motion estimation simple
        ARGS+=("-subq" "0")                  # Subpel quality minimum
        ARGS+=("-trellis" "0")               # Pas de trellis
        ARGS+=("-8x8dct" "0")                # Désactiver 8x8 DCT (Baseline)
        ARGS+=("-fast-pskip" "1")            # Fast P-skip
        ARGS+=("-wpredp" "0")                # Pas de weighted prediction
        ARGS+=("-mixed-refs" "0")            # Pas de mixed refs
        INJECT_DONE=1
    fi
done

# Exécuter FFmpeg avec les arguments modifiés
exec /usr/lib/jellyfin-ffmpeg/ffmpeg.original "${ARGS[@]}"
FFMPEG_WRAPPER

# Rendre le wrapper exécutable
chmod +x "$FFMPEG_WRAPPER"

echo "→ Wrapper FFmpeg installé ✅"
echo "→ Tous les transcodages H.264 utiliseront Baseline Profile !"
echo "==================================================="

